# findingsd.rc
#!/bin/sh
#

daemon="/usr/local/sbin/findingsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
